# templates/usr/lib/systemd/system/zeek-loggers.slice.j2
[Unit]
PartOf=zeek.target

[Slice]
MemoryMax={{ zeek_loggers_memory_max }}

[Install]
WantedBy=multi-user.target
