# templates/usr/lib/systemd/system/zeek-proxies.slice.j2
[Unit]
PartOf=zeek.target

[Slice]
MemoryMax={{ zeek_proxies_memory_max }}

[Install]
WantedBy=multi-user.target
