# templates/usr/lib/systemd/system/zeek.target.j2
[Unit]
Description=The Zeek Network Security Monitor

[Install]
WantedBy=multi-user.target

