# templates/usr/lib/systemd/system/zeek-workers.slice.j2
[Unit]
PartOf=zeek.target

[Slice]
MemoryMax={{ zeek_workers_memory_max }}

[Install]
WantedBy=multi-user.target
